# udev/98-bluez-pulse.rules
SUBSYSTEM=="input", GROUP="input", MODE="0660"
KERNEL=="input[0-9]*", RUN+="/opt/paradium/scripts/bluez-pulse-plugged"

